# schedule.cron
# NOTE: Change "USER" to your username

# Main Crontab for posting once a week.
0 18 * * 5 cd /home/USER && /usr/bin/python /home/USER/SOTW-bot RedditBot

# Test Crontab, will run script once every minute.
#* * * * * cd /home/USER && /usr/bin/python /home/USER/SOTW-bot RedditBotTest
