# schedule.cron
#Main Crontab for posting once a week.
0 18 * * 5 cd /home/jake && /usr/bin/python /home/jake/SOTW-bot

#Test Crontab, will run script once every minute.
#0/1 * * * * cd /home/jake && /usr/bin/python /home/jake/SOTW-bot
