# schedule.cron
0 18 * * 5 cd /home/jake && /usr/bin/python /home/jake/SOTW-bot
